# demo/pakage/arm64/DEBIAN/postinst
#!/bin/sh

echo "aibox installing"

systemctl enable yj-aibox
systemctl start yj-aibox

systemctl enable yj-mediaserver
systemctl start yj-mediaserver
